# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:redis]
command=/usr/bin/redis-server
autostart=true
autorestart=true

[program:caddy]
command=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true

[program:app]
command=python /app/main.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
